# smpcpd.service
[Unit]
Description=smpcp daemon
Documentation=man:smpcp(1)
After=mpd.service

[Service]
ExecStart=/usr/bin/smpcpd

[Install]
WantedBy=default.target
