# smpcpd.service
[Unit]
Description=smpcp daemon
Documentation=man:smpcp(1)

[Service]
ExecStart=/usr/bin/smpcpd

[Install]
WantedBy=multi-user.target
